# roles/redis/templates/root/etc/init/redis-server.conf.j2
description "redis server"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

setuid redis
setgid redis

exec /usr/bin/redis-server {{ redis.config_path }}
respawn
